# systemd/speedtest-cli.service
[Unit]
Description=Speedtest CLI Monitor - Run internet speed test

[Service]
Type=oneshot
WorkingDirectory=/home/je/SpeedtestMonitor/cli
ExecStart=/bin/bash /home/je/SpeedtestMonitor/cli/run-speedtest.sh 36738
StandardOutput=append:/home/je/SpeedtestMonitor/speedtest.log
StandardError=append:/home/je/SpeedtestMonitor/speedtest.log
